# ruby/supervisor.conf
[program:swarm-slave]
command=/usr/bin/swarm_slave.sh
user=root
